# infrastructure/artisan.crontab
* * * * * php-fpm php /app/artisan schedule:run > /tmp/stdout.log 2>&1
